# guestbook_sync.service
[Unit]
Description=Sync messages folder to Google Drive

[Service]
Type=oneshot
ExecStart=/opt/guestbook/sync_messages.sh
